# deploy/systemd/arklatex-watchdog.timer
[Unit]
Description=Run the ArkLaTex frozen-frame watchdog every 2 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
